# Redis-Cluster/7001/upstart.conf
description "redis server on port 7001"

start on runlevel [23]
stop on shutdown

exec sudo -u redis /usr/bin/redis-server /etc/redis/redis-7001.conf

respawn